# uviewd/crontab
# min hr mday mon wday
  */5 *  *    *   *	~/code/proj/uview/uviewd/runuview.sh
