# uviewd/crontab
# min hr mday mon wday
  */5 *  *    *   *	cd ~/code/proj/uview/uviewd/ && screen -S uview -d -m perl uviewd.pl
